# deploy/target/config/etc/systemd/system/restore-date.service
[Unit]
Description=Restore date from timestamp file

[Service]
Type=oneshot
RemainAfterExit=true
ExecStop=/bin/true
ExecStart=/usr/local/sbin/restore-date

[Install]
WantedBy=multi-user.target
