# deploy/target/config/etc/systemd/system/smith.service
[Unit]
Description=Smith firmware

[Service]
Type=forking
PIDFile=/var/run/smith.pid
ExecStart=/usr/local/sbin/start-smith
ExecStop=/usr/local/sbin/stop-smith

[Install]
WantedBy=multi-user.target
