# deploy/target/config/smith-common/etc/systemd/system/i2c-2.service
[Unit]
Description=Set pin mux to enable third I2C bus

[Service]
Type=oneshot
RemainAfterExit=true
ExecStop=/bin/true
ExecStart=/bin/bash -c '/bin/echo -n BB-I2C1 > /sys/devices/bone_capemgr.9/slots'

[Install]
WantedBy=multi-user.target
