# deploy/target/config/smith-common/etc/systemd/system/initialize-network.service
[Unit]
Description=Initialize network configuration via smith-config

[Service]
Type=oneshot
RemainAfterExit=true
Environment="GEM_HOME=/usr/local/lib/gems/1.9.1"
ExecStop=/bin/true
ExecStart=/usr/local/bin/smith-config init

[Install]
WantedBy=multi-user.target
