# deploy/target/config/smith-common/etc/systemd/system/smith-client.service
[Unit]
Description=Smith firmware web client
Wants=smith.service
After=smith.service

[Service]
Environment="GEM_HOME=/usr/local/lib/gems/1.9.1"
ExecStart=/usr/local/bin/smith-client --log-level=info --brief-log-format
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
