# deploy/target/config/smith-common/etc/systemd/system/smith-client.service
[Unit]
Description=Smith firmware web client

[Service]
Environment="GEM_HOME=/usr/local/lib/gems/1.9.1"
ExecStart=/usr/local/bin/smith-client
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
