# deploy/target/config/smith-common/etc/systemd/system/smith.service
[Unit]
Description=Smith firmware
Wants=owfs.service
After=owfs.service

[Service]
ExecStart=/usr/local/bin/smith

[Install]
WantedBy=multi-user.target
