# deploy/target/config/smith-common/etc/systemd/system/smith.service
[Unit]
Description=Smith firmware
Wants=owfs.service avr-firmware.service
After=owfs.service avr-firmware.service

[Service]
ExecStart=/usr/local/bin/smith --nostdio
Restart=on-failure
RestartSec=10
StartLimitInterval=110
StartLimitBurst=10
StartLimitAction=reboot-force

[Install]
WantedBy=multi-user.target
