# deploy/target/config/smith-development/etc/systemd/system/resize-rootfs.service
[Unit]
Description=Enlarge root filesystem partition to maximum capacity on first boot
ConditionPathExists=/var/local/resize-rootfs
After=remount-rootfs.service

[Service]
Type=oneshot
RemainAfterExit=true
ExecStartPost=/bin/rm /var/local/resize-rootfs
ExecStop=/bin/true
ExecStart=/sbin/resize2fs /dev/mmcblk0p2

[Install]
WantedBy=multi-user.target
